# deploy/systemd/zc.service
[Unit]
Description=zc local API and web workspace
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zc
Group=zc
WorkingDirectory=/opt/zc
EnvironmentFile=/etc/zc/zc.env
Environment=ENVIRONMENT=production
Environment=API_HOST=127.0.0.1
Environment=API_PORT=8000
Environment=API_WORKERS=1
Environment=STRICT_READINESS=true
Environment=UPLOAD_TEMP_DIR=/var/lib/zc/uploads
Environment=IDEMPOTENCY_DIR=/var/lib/zc/idempotency
Environment=CHAT_SESSION_DIR=/var/lib/zc/chat/sessions
ExecStart=/opt/zc/.venv/bin/zc --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
UMask=0077

StateDirectory=zc
StateDirectoryMode=0700
ReadWritePaths=/var/lib/zc
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
